# infra/packer/files/soctalk-firstboot.service
[Unit]
Description=SocTalk first-boot setup (k3s + helm install)
# Hook into cloud-init.target (not multi-user.target) to avoid the
# After=cloud-final.service ⇄ WantedBy=multi-user.target ordering
# cycle that systemd silently resolves by deleting our job.
After=cloud-init.target soctalk-setup-wizard.service network-online.target
Wants=network-online.target
ConditionPathExists=!/var/lib/soctalk-firstboot.done

[Service]
Type=oneshot
ExecStart=/usr/local/bin/soctalk-firstboot
RemainAfterExit=true
StandardOutput=journal+console
StandardError=journal+console
# No systemd start-timeout: this unit spends most of its time WAITING for the
# operator to finish the browser setup wizard (soctalk-setup-wizard.service is
# itself TimeoutStartSec=infinity — human-paced), then runs the k3s + helm
# install (helm has its own --wait --timeout). A fixed 20m start-timeout was
# SHORTER than firstboot.sh's own wait budget, so on slow/nested hardware
# systemd killed this unit mid-wait before the wizard was even submitted,
# leaving k3s up but SocTalk uninstalled. Let the script's internal budget +
# helm's timeout be the real gates.
TimeoutStartSec=infinity

[Install]
WantedBy=cloud-init.target
